# dlna-image/cfg/supervisor.conf
[supervisord]
nodaemon=true

[program:minidlna]
command=/etc/init.d/minidlna start

